# ansible/roles/kubelet/files/kubelet-override.service
[Service]
ExecStart=
ExecStart=/usr/bin/kubelet --config /etc/kubernetes/kubelet.yaml
